# overlay/etc/s6/gitaly/run
#!/bin/sh

echo "Starting gitaly.."
exec 2>&1
su-exec git gitaly /etc/gitlab/gitaly/config.toml
