# overlay/etc/s6/unicorn/run
#!/bin/sh

echo "Starting unicorn.."
exec 2>&1
cd /home/git/gitlab
rm -f /home/git/gitlab/tmp/pids/unicorn.pid
su-exec git bundle exec unicorn_rails -E production -c /etc/gitlab/unicorn.rb
