# overlay/etc/s6/nginx/run
#!/bin/sh

s6-svwait -t 5000 -u /etc/s6/unicorn
echo "Starting nginx.."
exec 2>&1
rm -f /tmp/nginx.pid
nginx -g "daemon off;pid /tmp/nginx.pid;"
