# overlay/etc/s6/sidekiq/run
#!/bin/sh

echo "Starting sidekiq.."
exec 2>&1
cd /home/git/gitlab
su-exec git bundle exec sidekiq -e production -C /home/git/gitlab/config/sidekiq_queues.yml
