# overlay/etc/s6/sshd/run
#!/bin/execlineb -P

if { echo "Starting sshd.." }
fdmove -c 2 1
/usr/sbin/sshd -D -f /etc/ssh/sshd_config -e
