# overlay/etc/s6/sshd/run
#!/bin/sh

echo "Starting sshd.."
exec 2>&1
/usr/sbin/sshd -D -f /etc/ssh/sshd_config -e
